# configs/supervisor.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --emperor /var/www/uwsgi

[program:nginx]
command = /usr/sbin/nginx